# dump/from_mtd/mtdblock2-extracted/ProgramFiles/etc/crontabs/root
*/2 * * * * /usr/bin/vihomed_monitor.sh
